# scripts/zeroclaw.service
[Unit]
Description=ZeroClaw AI Hardware Agent
Documentation=https://github.com/zeroclaw/zeroclaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@@RPI_USER@@
SupplementaryGroups=gpio spi i2c
WorkingDirectory=/home/@@RPI_USER@@/zeroclaw
ExecStart=/home/@@RPI_USER@@/zeroclaw/zeroclaw gateway start
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/@@RPI_USER@@/zeroclaw/.env
Environment=RUST_LOG=info

# Expand ~ in config path
Environment=HOME=/home/@@RPI_USER@@

[Install]
WantedBy=multi-user.target
